# leadnexconnect-v2/deploy/leadnex-web.service
[Unit]
Description=LeadNexConnect Web Frontend (Next.js)
After=network.target leadnex-api.service
Wants=leadnex-api.service

[Service]
Type=simple
User=abu-lukas
Group=abu-lukas

WorkingDirectory=/home/abu-lukas/LeadNexConnect/leadnexconnect-v2/apps/web

# Load environment variables
EnvironmentFile=/home/abu-lukas/LeadNexConnect/leadnexconnect-v2/apps/web/.env.local

# Next.js environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0

# Start Next.js production server
ExecStart=/home/abu-lukas/LeadNexConnect/leadnexconnect-v2/node_modules/.bin/next start

# Gracefully stop
KillSignal=SIGTERM
TimeoutStopSec=30
KillMode=mixed

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging — view with:  journalctl -u leadnex-web -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=leadnex-web

# Basic security hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
