# leadnexconnect-v2/deploy/leadnex-api.service
[Unit]
Description=LeadNexConnect API Server (Express/Node.js)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=abu-lukas
Group=abu-lukas

WorkingDirectory=/home/abu-lukas/LeadNexConnect/leadnexconnect-v2/apps/api

# Load environment variables from the .env file
EnvironmentFile=/home/abu-lukas/LeadNexConnect/leadnexconnect-v2/apps/api/.env

# Always force production mode
Environment=NODE_ENV=production

# Start the compiled API bundle
ExecStart=/usr/bin/node /home/abu-lukas/LeadNexConnect/leadnexconnect-v2/apps/api/dist/index.js

# Gracefully stop
KillSignal=SIGTERM
TimeoutStopSec=30
KillMode=mixed

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging — view with:  journalctl -u leadnex-api -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=leadnex-api

# Basic security hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
